# tools/launchagents/io.orchestrateos.jarvis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.orchestrateos.jarvis</string>
    <key>WorkingDirectory</key>
    <string>/Users/srinivas/Library/Application Support/OrchestrateOS</string>
    <key>ProgramArguments</key>
    <array>
        <string>python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>jarvis:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>5004</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/srinivas/Library/Logs/OrchestrateOS/jarvis.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/srinivas/Library/Logs/OrchestrateOS/jarvis_error.log</string>
</dict>
</plist>
